# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn --timeout 120 -w 4 -b 0.0.0.0:8000 wsgi:app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=python worker.py
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
